# snowblocks/xdg/autostart/dropbox.desktop
[Desktop Entry]
Type=Application
Name=Dropbox
Exec=dropbox
Icon=dropbox
